# temp.crontab
*/1 * * * * ~/sync/sync.sh  &> ~/sync/sync.out 
